# tge/systemd/tge-apply.path
[Unit]
Description=TGE: Trigger apply when tun0 appears

[Path]
PathExists=/sys/class/net/tun0
Unit=tge-apply.service

[Install]
WantedBy=multi-user.target